# dp-build/services/play/run
#!/bin/sh
exec 2>&1
cd /usr/dp-app/bin && exec ./dp-app -Dconfig.resource=../conf/docker.conf
